# systemd/car.service.j2
[Unit]
Description=Raspberry PI Car service

[Service]
User={{ service_user }}
Group={{ service_user }}
Restart=on-failure
WorkingDirectory={{ bin_dir }}
ExecStart={{ bin_dir }}/{{ bin_file }}
